# bin/adc-key/adckeys.service
[Unit]
Description=External service

[Service]
ExecStart=/usr/local/bin/adckeys.sh
Restart=always
Type=simple
RestartSec=0.5

[Install]
WantedBy=multi-user.target